# delivery-kid/ansible/roles/pinning-service/templates/sync-releases.service.j2
[Unit]
Description=Sync releases from PickiPedia to delivery-kid
After=network.target pinning-service.service ipfs.service
Wants=pinning-service.service

[Service]
Type=oneshot
WorkingDirectory=/opt/pinning-service
ExecStart=/usr/bin/node sync-releases.js --apply
Environment=WIKI_API=https://pickipedia.xyz/api.php
Environment=PINNING_API=http://localhost:3001
Environment=REQUIRED_GROUP=release

# Give IPFS time to fetch content
TimeoutStartSec=600
