# delivery-kid/ansible/roles/pinning-service/templates/sync-releases.timer.j2
[Unit]
Description=Sync releases from PickiPedia hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
